# iRony/debian/preinst
#!/bin/sh

set -e

reload_apache()
{
    if apache2ctl configtest 2>/dev/null; then
        invoke-rc.d apache2 $1 || true
    else
        echo "Your apache2 configuration is broken, so we're not restarting it for you."
    fi
}


case "$1" in
    upgrade|install)
        # rename old apache/http sites rom sitename to sitename.conf
        if [ -e /etc/apache2/sites-available/iRony -a ! -e /etc/apache2/sites-available/iRony.conf ]; then
            if [ -h /etc/apache2/sites-enabled/iRony ]; then
                a2dissite iRony
            fi
            mv /etc/apache2/sites-available/iRony /etc/apache2/sites-available/iRony.conf
            a2ensite iRony.conf
            reload_apache
        fi
    ;;
esac

exit 0

