# pykolab/debian/kolab-saslauthd.postinst
#!/bin/bash

#DEBHELPER#

set -e

case "$1" in
    configure|upgrade)
        invoke-rc.d kolab-saslauthd restart || true
    ;;

esac
