# pykolab/debian/wallace.postinst
#!/bin/bash

#DEBHELPER#

set -e

case "$1" in
    install|configure)
        mkdir -p /var/spool/pykolab/wallace && chown -R kolab:kolab-n /var/spool/pykolab
    ;;

    upgrade)
        mkdir -p /var/spool/pykolab/wallace && chown -R kolab:kolab-n /var/spool/pykolab
        invoke-rc.d wallace restart || true
    ;;

esac
